# chain2gate.service
[Unit]
Description=Chain2gate

[Service]
ExecStart=/home/pi/chain2gate/chain2gate_client.py
Restart=always

[Install]
WantedBy=multi-user.target
